# deploy/systemd/crypto-radar-backup.service
[Unit]
Description=Crypto Radar database backup

[Service]
Type=oneshot
WorkingDirectory=/opt/crypto_radar
ExecStart=/opt/crypto_radar/venv/bin/python /opt/crypto_radar/backup_database.py
User=root